# agents/agents/systemd/dci-agent-dci.timer
[Unit]
Descriptuon=The DCI dci agent timer

[Timer]
OnCalendar=00:20:00
Persistent=true

[Install]
WantedBy=timers.target
