# feeders/systemd/dci-feeder-dci.service
[Unit]
Description=The DCI dci feeder

[Service]
Type=oneshot
ExecStart=dci-feeder-dci
