# feeders/systemd/dci-feeder-github.service
[Unit]
Description=The DCI github feeder

[Service]
Type=oneshot
ExecStart=dci-feeder-github https://github.com/%I
