# feeders/systemd/dci-feeder-dci.timer
[Unit]
Descriptuon=The DCI dci feeder timer

[Timer]
OnCalendar=00:10:00
Persistent=true

[Install]
WantedBy=timers.target
